# deploy/display_setup.service
[Unit]
Description=Set HDMI display resolution and rotation
After=graphical.target

[Service]
User=admin
Environment=DISPLAY=:0
ExecStart=/home/admin/set-display.sh
Restart=on-failure

[Install]
WantedBy=graphical.target